# s6-overlay/services.d/reverb/run
#!/usr/bin/with-contenv sh


cat <<EOL
┌────────────────────┐
       REVERB
└────────────────────┘
EOL

cd /var/www/html
exec s6-setuidgid nginx php artisan reverb:start --debug

