# s6-overlay/services.d/php-fpm/run
#!/usr/bin/with-contenv bash


cat <<EOL
┌──────────┐
  PHP-FPM
└──────────┘
EOL


exec php-fpm -F
